# roles/init/templates/networkcheck.timer.j2
[Unit]
Description=Network check Timer

[Timer]
OnUnitActiveSec=600s
OnBootSec=600s

[Install]
WantedBy=timers.target
